# roles/uberspace_postgresql/templates/postgresql.ini.j2
[program:postgresql]
command=postgres -D %(ENV_HOME)s/{{ postgres_projectsdir }}{{ postgres_dir_name }}/data/
autostart=yes
autorestart=yes
startsecs=15
