# roles/uberspace_prosody/templates/prosody.ini.j2
[program:prosody]
command=/bin/bash -c "source {{ uberspace_spacedir }}/.bash_profile && prosody"
autostart=yes
autorestart=yes
startretries=1
stopasgroup=true
killasgroup=true
stopsignal=INT
