# roles/uberspace_coturn/templates/coturn.ini.j2
[program:coturn]
command={{ coturn_install_dir }}/bin/turnserver -c {{ coturn_configfile }}
